# docker/root/etc/services.d/cnl2feedcrawler/run
#!/usr/bin/with-contenv bash

exec cnl2feedcrawler --url=$URL
